# deploy/com.eyeguard.vault.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- EyeGuard vault daemon — the privileged half of the locked-down build.
     Installed to /Library/LaunchDaemons (root-owned) by setup-lockdown.sh.
     Runs as ROOT: holds the key, does all network writes. The monitored user
     cannot read, edit, or unload it (that needs admin, which only Dad has). -->
<plist version="1.0"><dict>
  <key>Label</key><string>com.eyeguard.vault</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/EyeGuard.app/Contents/Resources/python/bin/python3.12</string>
    <string>-m</string>
    <string>eyeguard.vault</string>
    <string>--config</string>
    <string>/Library/Application Support/EyeGuard/config.yaml</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Library/Application Support/EyeGuard</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key><string>/Library/Application Support/EyeGuard</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Library/Application Support/EyeGuard/logs/vault.out.log</string>
  <key>StandardErrorPath</key><string>/Library/Application Support/EyeGuard/logs/vault.err.log</string>
</dict></plist>
